# scripts/aqu.service
[Unit]
Description=aqu VPN panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/aqu serve
Restart=always
RestartSec=3
User=root
# позволяет панели биндить привилегированные порты, если понадобится
AmbientCapabilities=CAP_NET_BIND_SERVICE
# каталог с конфигом и данными
Environment=AQU_DIR=/etc/aqu

[Install]
WantedBy=multi-user.target
